# opencharacters.service
[Unit]
Description=OpenCharacters Web Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/opencharacterss
Environment=PORT=4000
Environment=HOST=0.0.0.0
Environment=PUBLIC_IP=100.87.19.163
Environment=PUBLIC_URL=http://100.87.19.163:4000/
Environment=IFACE=eth0
ExecStart=/opt/opencharacterss/start-service.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
